# plugin/src/Install.service
[Unit]
Description="install util disable hotkeys"
 
[Service]
ExecStart=/usr/bin/install.sh
 
[Install]
WantedBy=multi-user.target